# plexmediaserver.conf
[program:plexmediaserver]
user = nobody
command = /home/nobody/setup.sh
