# rar2fs.conf
[program:rar2fs]
user = nobody
command = /home/nobody/rar2fs.sh