# build/plexmediaserver.conf
[program:start-script]
user = nobody
command = /home/nobody/start.sh
umask = 000

[program:claim-script]
user = nobody
command = /home/nobody/claim.sh
umask = 000
autorestart = false
startsecs = 0