# build/plexmediaserver.conf
[program:plexmediaserver]
user = nobody
command = /home/nobody/start.sh
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh 'Plex Media Server'
umask = 000
