# setup/plexmediaserver.conf
[program:plexmediaserver]
user = nobody
command = /home/nobody/start.sh
